# packages/luci-app-remotbot/files/etc/uci-defaults/luci-app-remotbot
#!/bin/sh
uci -q get remotbot.main >/dev/null 2>&1 || uci set remotbot.main=remotbot
uci -q get remotbot.main.bot_token >/dev/null 2>&1 || uci set remotbot.main.bot_token=''
uci -q get remotbot.main.allowed_users >/dev/null 2>&1 || uci set remotbot.main.allowed_users=''
uci -q get remotbot.main.cgi_online_path >/dev/null 2>&1 || uci set remotbot.main.cgi_online_path='/www/cgi-bin/online'
uci -q get remotbot.main.enabled >/dev/null 2>&1 || uci set remotbot.main.enabled='0'
uci -q get remotbot.main.log_level >/dev/null 2>&1 || uci set remotbot.main.log_level='INFO'
uci commit remotbot 2>/dev/null || true
exit 0
