# packages/luci-app-remotbot/files/etc/uci-defaults/luci-app-remotbot
#!/bin/sh
if ! uci -q get remotbot.main >/dev/null 2>&1; then
    uci set remotbot.main=remotbot
    uci set remotbot.main.bot_token=''
    uci set remotbot.main.allowed_users=''
    uci set remotbot.main.cgi_online_path='/www/cgi-bin/online'
    uci set remotbot.main.enabled='0'
    uci set remotbot.main.log_level='INFO'
    uci commit remotbot
fi
exit 0
