# luci-app-remotbot/root/etc/init.d/remotbot
#!/bin/sh /etc/rc.common
START=99; STOP=10; USE_PROCD=1
PROG=/usr/bin/python3
SCRIPT=/usr/share/remotbot/pi4Bot.py
PIDFILE=/var/run/remotbot.pid
CONFIG=/etc/remotbot/config.json
start_service() {
    [ ! -f "$CONFIG" ] && echo "Config not found: $CONFIG" && return 1
    TOKEN=$(jsonfilter -i "$CONFIG" -e '@.bot_token' 2>/dev/null)
    [ -z "$TOKEN" ] && echo "bot_token not configured" && return 1
    procd_open_instance
    procd_set_param command $PROG $SCRIPT
    procd_set_param pidfile $PIDFILE
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param user root
    procd_close_instance
}
stop_service() {
    [ -f "$PIDFILE" ] && PID=$(cat $PIDFILE) && kill -TERM $PID 2>/dev/null && sleep 2 && kill -9 $PID 2>/dev/null
    rm -f $PIDFILE
}
